# Misc/payload-startup.service
[Unit]
Description=Rocketry@VT Payload Startup Service


[Service]
Type=simple
#if service finishes [not including manually stopping the service], restart after 10 sec
Restart=always
RestartSec=10
ExecStart=/home/debian/Payload/payload.py

[Install]
WantedBy=default.target
